# package/luci/themes/luci-theme-witi-material/root/etc/uci-defaults/30_luci-theme-material
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Material=/luci-static/witi-material
	set luci.main.mediaurlbase=/luci-static/witi-material
	commit luci
EOF
exit 0
